# noinvade-webapp.service
[Unit]
Description=noinvade java web app
After=network.target

[Service]
ExecStart=/home/ubuntu/noinvade/start.sh

[Install]
WantedBy=default.target